# fcgiwrap/supervisord.conf
[program:fcgiwrap]
command=/usr/bin/fcgiwrap -f -s unix:/var/run/fcgiwrap.sock
