# mailman/supervisord.conf
[program:mailman]
command=/usr/local/bin/foreground -c '/usr/local/mailman/bin/mailmanctl -s start' -p /data/mailman/data/master-qrunner.pid
