# mailman/supervisord.conf
[program:mailman]
command=/usr/local/bin/mailman.sh
