# panel/supervisord.conf
[program:panel]
command=/usr/bin/spawn-fcgi -u www-data -g www-data -s /var/run/panel.sock -n -- /usr/local/bin/panel
