# postsrsd/supervisord.conf
[program:postsrsd]
command=/usr/sbin/postsrsd -s/etc/postsrsd.secret -dpfohl.space -Xmars.pfohl.space -unobody -c/var/lib/postsrsd
