# postsrsd/supervisord.conf
[program:postsrsd]
command=/usr/sbin/postsrsd -s/etc/postsrsd.secret -dpfohl.email -Xlists.pfohl.email -unobody -c/var/lib/postsrsd
