# postsrsd/supervisord.conf
[program:postsrsd]
command=/usr/sbin/postsrsd -s/etc/postsrsd.secret -dexample.com -Xlists.example.com -unobody -c/var/lib/postsrsd
